# systemd/servedash.socket
[Unit]
Description=ServeDash Socket

[Socket]
ListenStream=7098
Accept=no
NoDelay=true

[Install]
WantedBy=sockets.target